# templates/supervisor.ini.j2
# This file was generated by Ansible for {{ ansible_fqdn }}
# Do NOT modify this file by hand!

[program:{{wsgi_supervisor_process_name}}]
autorestart=true
autostart=true
command={{wsgi_virtualenv}}/bin/uwsgi --ini {{wsgi_cfg_dir}}/uwsgi.ini
redirect_stderr=true
stdout_logfile={{wsgi_log_dir}}/uwsgi.log
user={{wsgi_user}}

{% if wsgi_supervisor_httpok %}
[eventlistener:httpok-{{wsgi_supervisor_process_name}}]
command={{wsgi_supervisor_httpok}}
events=TICK_60
{% endif %}
