# templates/supervisor.ini.j2
# This file was generated by Ansible for {{ ansible_fqdn }}
# Do NOT modify this file by hand!

[program:uwsgi-{{ base_project_name }}]
autorestart=true
autostart=true
command={{ wsgi_virtualenv }}/bin/uwsgi --ini {{ base_configuration_directory }}/uwsgi.ini
redirect_stderr=true
stdout_logfile={{ base_logs_directory }}/uwsgi.log
user={{ base_deploy_user }}
