# templates/upstart-uwsgi.conf.j2
#!upstart

# This file was generated by Ansible for {{ansible_fqdn}}
# Do NOT modify this file by hand!

{% set app = item -%}
{% set run_dir = item.run_dir|default("/opt/%s/run" % item.name) -%}
{% set log_dir = item.log_dir|default("/opt/%s/log" % item.name) -%}
{% set etc_dir = item.etc_dir|default("/opt/%s/etc" % item.name) -%}
{% set env_dir = item.virtualenv|default("/opt/%s/env" % item.name) -%}
{% set app_dir = item.app_dir|default("/opt/%s/src" % item.name) -%}

description "Start uwsgi for {{app.name}}"

start on (local-filesystems and runlevel [2345])
stop on runlevel [!2345]

{% for key in app.env|default({}) %}
env {{key}}={{app.env[key]}}
{% endfor %}

# expect fork
respawn

limit nofile 65536 65536
console log
chdir {{app_dir}}
setuid {{wsgi_user}}
setgid {{wsgi_group}}

exec {{env_dir}}/bin/uwsgi --ini {{etc_dir}}/{{app.name}}-uwsgi.ini
