# templates/upstart-uwsgi.conf.j2
#!upstart

# This file was generated by Ansible for {{ansible_fqdn}}
# Do NOT modify this file by hand!

{% set app = item -%}

description "Start uwsgi for {{app.name}}"

start on (local-filesystems and runlevel [2345])
stop on runlevel [!2345]

{% for key in app.env|default({}) %}
env {{key}}={{app.env[key]}}
{% endfor %}

# expect fork
respawn

limit nofile 65536 65536
console log
chdir {{app.app_dir|default(wsgi_app_dir)}}
setuid {{wsgi_user}}
setgid {{wsgi_group}}

exec {{app.virtualenv|default(wsgi_virtualenv)}}/bin/uwsgi --ini {{app.etc_dir|default(wsgi_etc_dir)}}/{{app.name}}-uwsgi.ini
