# templates/upstart-gunicorn.conf.j2
#!upstart

# This file was generated by Ansible for {{ansible_fqdn}}
# Do NOT modify this file by hand!

{% set app = item -%}

description "Start Gunicorn server for {{app.name}}"

start on (filesystem)
stop on runlevel [016]

# expect fork
respawn

setuid {{wsgi_user}}
setgid {{wsgi_group}}

limit nofile 65536 65536
console log
chdir {{app.app_dir|default(wsgi_app_dir)}}

exec {{app.virtualenv|default(wsgi_virtualenv)}}/bin/{{app.gunicorn_command|default(wsgi_gunicorn_command)}} \
    -n {{app.name}} \
    -p {{app.run_dir|default(wsgi_run_dir)}}/{{app.name}}.pid \
    -k {{app.gunicorn_worker|default(wsgi_gunicorn_worker)}} \
    -w {{app.gunicorn_workers|default(wsgi_gunicorn_workers)}} \
    -b {{app.socket|default('unix:' + app.get('run_dir', wsgi_run_dir) + '/' + app.name)}} \
    --log-file {{app.log_dir|default(wsgi_log_dir)}}/{{app.name}}-wsgi.log \
    {% if wsgi_reload %}--reload {% endif %}{{app.module|default('wsgi:app')}}
